# files/jenkins-agent.RedHat
#!/bin/sh
#
# /etc/init.d/jenkins-agent
# init script for Jenkins Agent Swarm Client.
#
# chkconfig: 2345 90 60
# description: Jenkins Agent Swarm Client
#
RETVAL=0

NAME=jenkins-agent
JENKINS_CONFIG=/etc/sysconfig/$NAME
LOCK_FILE=/var/lock/subsys/$NAME

# Source function library.
. /etc/init.d/functions

# Read config
[ -f "$JENKINS_CONFIG" ] && . "$JENKINS_CONFIG"

if [ -x /sbin/runuser ] ; then
    RUNUSER=runuser
else
    RUNUSER=su
fi

agent_start() {
  echo Starting Jenkins Agent...

  # the default location is /var/run/jenkins/jenkins.pid but the parent directory needs to be created
  mkdir `dirname $PIDFILE` > /dev/null 2>&1 || true
  chown $JENKINS_AGENT_USER `dirname $PIDFILE`

  # create log directory
  mkdir -p `dirname $JENKINS_AGENT_LOG` > /dev/null 2>&1 || true
  chown $JENKINS_AGENT_USER -R `dirname $JENKINS_AGENT_LOG`

  # --user in daemon doesn't prepare environment variables like HOME, USER, LOGNAME or USERNAME,
  # so we let su do so for us now
  $RUNUSER - $JENKINS_AGENT_USER -c "$JAVA $JAVA_ARGS -jar $JENKINS_AGENT_JAR $JENKINS_AGENT_ARGS >> $JENKINS_AGENT_LOG 2>&1 &"
  pgrep -f -u $JENKINS_AGENT_USER $JENKINS_AGENT_JAR > $PIDFILE
  RETVAL=$?
  [ $RETVAL -eq 0 ] && touch $LOCK_FILE
}
agent_stop() {
  echo Stopping Jenkins Agent...
  pid=`cat $PIDFILE`

  killproc -p $PIDFILE

  # Wait until the monitor exits
  while (checkpid $pid)
  do
      usleep 100000
  done

  RETVAL=$?
  [ $RETVAL -eq 0 ] && rm -f $LOCK_FILE
}

agent_restart() {
  echo Restarting Jenkins Agent...
  agent_stop
  agent_start

  RETVAL=$?
  [ $RETVAL -eq 0 ] && touch $LOCK_FILE
}
agent_status() {
  echo Jenkins Agent status:
  status -p $PIDFILE
  RETVAL=$?
}
case "$1" in
  start)
    agent_start
    ;;
  stop)
    agent_stop
    ;;
  restart)
    agent_restart
    ;;
  status)
    agent_status
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac

exit $RETVAL
